# units/web.tags.expire.timer
[Unit]
Description=Expire Donation Tags
Requires=docker.service

[Timer]
OnCalendar=*:0/5

[X-Fleet]
X-ConditionMachineOf=web.tags.expire.service
